# etc/anyorigin.service.in
[Unit]
Description=anyorigin
Wants=nginx

[Service]
ExecStart={{PREFIX}}/usr/bin/anyorigin --address {{HOST}}:{{PORT}}

[Install]
WantedBy=multi-user.target
